# crontab
#run python script every 4 hour

0 */4 * * * /usr/local/bin/python /app/pynod/update.py > /proc/1/fd/1 2>/proc/1/fd/2
